# roles/go/templates/etc/profile.d/go-bin.sh.j2
{{ ansible_managed | comment }}

export PATH=$PATH:/usr/local/go/bin
